# scripts/systemd/docker-iptables.service
[Unit]
Description=Docker DOCKER-USER firewall rules (allow container outbound)
After=docker.service
Requires=docker.service
BindsTo=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Rebuild DOCKER-USER on every docker start/reboot – must include all compose bridge subnets.
ExecStart=/usr/sbin/iptables -F DOCKER-USER
ExecStart=/usr/sbin/iptables -I DOCKER-USER 1 -m conntrack --ctstate RELATED,ESTABLISHED -j RETURN
ExecStart=/usr/sbin/iptables -I DOCKER-USER 2 -s 127.0.0.0/8 -j RETURN
ExecStart=/usr/sbin/iptables -I DOCKER-USER 3 -s 192.168.11.0/24 -j RETURN
ExecStart=/usr/sbin/iptables -I DOCKER-USER 4 -s 172.17.0.0/16 -j RETURN
ExecStart=/usr/sbin/iptables -I DOCKER-USER 5 -s 172.18.0.0/16 -j RETURN
ExecStart=/usr/sbin/iptables -I DOCKER-USER 6 -s 172.19.0.0/16 -j RETURN
ExecStart=/usr/sbin/iptables -I DOCKER-USER 7 -s 172.20.0.0/16 -j RETURN
ExecStart=/usr/sbin/iptables -I DOCKER-USER 8 -s 172.22.0.0/16 -j RETURN
ExecStart=/usr/sbin/iptables -I DOCKER-USER 9 -s 172.23.0.0/16 -j RETURN
ExecStart=/usr/sbin/iptables -A DOCKER-USER -j DROP

[Install]
WantedBy=multi-user.target
